# Configs/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate desktop wallpaper (archScratch)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-rotate
